# docs/kmscon-login@.service.in
[Unit]
Description=KMS System Console login program on /dev/pts/%I
Documentation=man:kmscon(1)
ConditionPathExists=/dev/pts/%I

[Service]
ExecStart=/sbin/agetty -o '-p -- \\u' --noclear pts/%I -
Restart=always
RestartSec=0
TTYReset=yes
TTYVHangup=yes
IgnoreSIGPIPE=no
SendSIGHUP=yes
