# docs/kmscon.service.in
[Unit]
Description=KMS System Console
Documentation=man:kmscon(1)
After=plymouth-quit-wait.service
After=systemd-user-sessions.service
After=rc-local.service

[Service]
ExecStart=@bindir@/kmscon --login -- /usr/bin/systemctl start kmscon-login@{pts_number}

[Install]
WantedBy=multi-user.target
